# Webmin-y-Virtualmin--7a5b69265dabb75e98929f6bcaf64b5924996415/auto-tunnel.service
[Unit]
Description=Sistema de Túnel Automático 24/7
After=network.target network-online.target
Wants=network-online.target
Requires=network.target

[Service]
Type=simple
ExecStart=/bin/bash /root/auto_tunnel_system.sh start
ExecStop=/bin/bash /root/auto_tunnel_system.sh stop
ExecReload=/bin/bash /root/auto_tunnel_system.sh restart
Restart=always
RestartSec=10
User=root
Group=root
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
StandardOutput=journal
StandardError=journal
SyslogIdentifier=auto-tunnel

# Configuración de seguridad
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/etc/auto_tunnel_config.conf /var/log/auto_tunnel_system.log /var/run
ProtectHome=yes
PrivateDevices=yes

# Límites de recursos
LimitNOFILE=65536
MemoryLimit=256M
CPUQuota=50%

[Install]
WantedBy=multi-user.target