# Webmin-y-Virtualmin--7a5b69265dabb75e98929f6bcaf64b5924996415/ssl_monitor.timer
[Unit]
Description=Temporizador para Monitoreo SSL cada hora
Requires=ssl_monitor.service

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target