# Webmin-y-Virtualmin--7a5b69265dabb75e98929f6bcaf64b5924996415/ssl_monitor.service
[Unit]
Description=Sistema de Monitoreo SSL Avanzado
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/advanced_ssl_manager.sh validate
StandardOutput=journal
StandardError=journal