# auto-tunnel.service
[Unit]
Description=Disabled insecure auto tunnel service

[Service]
Type=oneshot
ExecStart=/bin/false

[Install]
WantedBy=multi-user.target
