# com.webmin.repovalidation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.webmin.repovalidation</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/webmin-repo-validation/webmin-repo-validation.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>/var/log/webmin-repo-validation.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/webmin-repo-validation-error.log</string>
    <key>WorkingDirectory</key>
    <string>/opt/webmin-repo-validation</string>
</dict>
</plist>
