# webmin-performance-optimizer.service
[Unit]
Description=Sistema de Optimización Ultra-Alta Performance Webmin
After=network.target apache2.service mysql.service
Requires=network.target

[Service]
Type=oneshot
ExecStart=/opt/webmin-performance/webmin-performance-optimizer.sh optimize
ExecStartPost=/opt/webmin-performance/webmin-performance-optimizer.sh monitor
RemainAfterExit=true
User=root
WorkingDirectory=/opt/webmin-performance

[Install]
WantedBy=multi-user.target
