# webmin-self-healing.service
[Unit]
Description=Sistema de Auto-Reparación Inteligente Webmin/Virtualmin
After=network.target

[Service]
Type=simple
ExecStart=/opt/webmin-self-healing/auto-repair.sh
Restart=always
RestartSec=10
User=root

[Install]
WantedBy=multi-user.target
