# Webmin-y-Virtualmin--7a5b69265dabb75e98929f6bcaf64b5924996415/cluster_infrastructure/ansible/roles/cost_monitoring/templates/cost_optimization.timer.j2
[Unit]
Description=Timer for AWS Cost Optimization - Daily execution
Requires=cost-optimization.service

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target